# nexus.supervisor.conf
[program:nexus]
command=/u01/nexus-2.9.0-04/bin/nexus console
user=www-data
